# specs/ext-sources/php-fpm72jposug.xml
<?xml version="1.0"?>
<!-- This is based on https://gist.github.com/webjay/888451 -->
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="php-fpm72jposug">
  <service name="network/php-fpm72jposug" type="service" version="1">
    <create_default_instance enabled='false'/>
    <single_instance/>
    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default"/>
    </dependency>
    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local"/>
    </dependency>
    <dependent name="nginx" restart_on="none" grouping="optional_all">
      <service_fmri value="svc:/milestone/multi-user"/>
    </dependent>
    <method_context project=":default" resource_pool=":default" working_directory=":default"/>
    <exec_method type="method" name="start" exec="/opt/jposug/php/7.2/sbin/php-fpm -c /etc/php/7.2jposug/php.ini -y %{config_file}" timeout_seconds="60"/>
    <exec_method type="method" name="stop" exec=":kill" timeout_seconds="60"/>
    <exec_method type="method" name="refresh" exec=":kill -HUP" timeout_seconds="60"/>
    <property_group name="application" type="application">
      <propval name="config_file" type="astring" value="/etc/php/7.2jposug/php-fpm.conf"/>
    </property_group>
    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="contract"/>
      <propval name="ignore_error" type="astring" value="core,signal"/>
    </property_group>
    <stability value="Evolving"/>
    <template>
      <common_name>
	<loctext xml:lang="C">PHP FastCGI Process Manager</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
